# jarvis-gateway.service
[Unit]
Description=Jarvis Telegram Gateway
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/path/to/jarvis-telegram-gateway
ExecStart=/usr/bin/python3 /path/to/jarvis-telegram-gateway/gateway.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
